# restic/lunos/systemd/backup-reminder.timer
[Unit]
Description=Hourly restic backup notification check

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
